# inn-cnfsstat.cron
#!/bin/sh

. /usr/share/news/innshellvars

DOCNFSSTAT=`echo ${DOCNFSSTAT} | tr A-Z a-z`

if [ -z "${DOCNFSSTAT}" -o "${DOCNFSSTAT}" = "on" -o "${DOCNFSSTAT}" = "true" -o "${DOCNFSSTAT}" = "yes" ]; then
/usr/bin/cnfsstat -s
fi
