# inn-cnfsstat.cron
#!/bin/sh

. /usr/share/news/innshellvars
. /etc/rc.d/init.d/functions

if is_yes $DOCNFSSTAT ; then
    /usr/bin/cnfsstat -s
fi
