# inn-cnfsstat.cron
#!/bin/sh

. /usr/share/news/innshellvars

if `echo $DOCNFSSTAT | grep -q -e "^[tTyY].*"`; then
    /usr/bin/cnfsstat -s
fi
