# inn.crontab
MAILTO=""

01 * * * *	news	/usr/bin/nntpsend
01/10 * * * * 	news	/usr/bin/cnfsstat.cron
02 1 * * * 	news	/usr/bin/news.daily delayrm
02 1 * * * 	root	/usr/bin/rnews -U
